# .config/autostart/mathworks-service-host.desktop
[Desktop Entry]
Type=Application
Name=Mathworks Service Host
Exec=/home/quang/.MathWorks/ServiceHost/-mw_shared_installs/v2026.5.0.3/bin/glnxa64/MathWorksServiceHost service --realm-id companion@prod@production
Terminal=false
